# deploy/systemd/chat-grid.service
[Unit]
Description=Chat Grid signaling server
After=network.target

[Service]
Type=simple
User=bestmidi
Group=bestmidi
WorkingDirectory=/home/bestmidi/chgrid/server
Environment=PATH=/home/bestmidi/chgrid/server/.venv/bin:/usr/bin:/bin
ExecStartPre=/usr/bin/mkdir -p /home/bestmidi/chgrid/server/runtime
ExecStart=/home/bestmidi/chgrid/server/.venv/bin/python main.py --config /home/bestmidi/chgrid/server/config.toml
StandardOutput=append:/home/bestmidi/chgrid/server/runtime/server.log
StandardError=append:/home/bestmidi/chgrid/server/runtime/server.log
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
